# etc/crontabs/root
# 定时配置
# 格式：分 时 日 月 周 命令

# 每天5点重启
0 5 * * * /data/adb/modules/SAM/scripts/service.sh restart

# 每周一8点更新 GitHub 规则
0 8 * * 1 /data/adb/modules/SAM/scripts/host.sh gh
